# com.mmovsesyan.ollama-telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mmovsesyan.ollama-telegram-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>cd /Users/%USER%/ollama__telegram_bot && set -a && source .env && set +a && poetry run python main.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/%USER%/ollama__telegram_bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/%USER%/ollama__telegram_bot/bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/%USER%/ollama__telegram_bot/bot.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/%USER%/.local/bin</string>
    </dict>
</dict>
</plist>
